# exchange/etc/supervisor/conf.d/taler-exchange-httpd.conf
[program:taler-exchange-httpd]
command=/usr/bin/taler-exchange-httpd
user=taler-exchange-httpd
autostart=true
autorestart=true
