# merchant/etc/supervisor/conf.d/taler-merchant.conf
[program:taler-merchant]
command=/usr/bin/taler-merchant
user=taler-merchant-httpd
autostart=true
autorestart=true
